# hdwx.target
[Unit]
Description=HDWX Target
Requires=network-online.target
After=network-online.target

[Install]
WantedBy=multi-user.target